# dot_config/systemd/user/eclair.target
[Unit]
Description=Eclair Notification Center
PartOf=graphical-session.target

[Install]
WantedBy=graphical-session.target