# dot_config/systemd/user/eclair-ui.service
[Unit]
Description=Eclair GTK4 Notification UI
After=graphical-session.target
PartOf=eclair.target

[Service]
Type=simple
ExecStart=/home/bldnwine/owlcat/eclair/ui-rs/target/release/eclair-ui daemon
Restart=on-failure
RestartSec=5
Environment=WAYLAND_DISPLAY=wayland-1
KillMode=mixed

[Install]
WantedBy=eclair.target
