# dot_config/systemd/user/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (v2026.1.30)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart="/usr/bin/node" "/home/bldnwine/.npm-global/lib/node_modules/openclaw/dist/index.js" gateway --port 18789
Restart=always
RestartSec=5
KillMode=process
Environment=HOME=/home/bldnwine
Environment="PATH=/home/bldnwine/.local/bin:/home/bldnwine/.npm-global/bin:/home/bldnwine/bin:/home/bldnwine/.nvm/current/bin:/home/bldnwine/.fnm/current/bin:/home/bldnwine/.volta/bin:/home/bldnwine/.asdf/shims:/home/bldnwine/.local/share/pnpm:/home/bldnwine/.bun/bin:/usr/local/bin:/usr/bin:/bin"
Environment=OPENCLAW_GATEWAY_PORT=18789
Environment=OPENCLAW_GATEWAY_TOKEN=b13431e92faf2f04ed1d7cbfd8ee71b5832640a641aa8732
Environment="OPENCLAW_SYSTEMD_UNIT=openclaw-gateway.service"
Environment=OPENCLAW_SERVICE_MARKER=openclaw
Environment=OPENCLAW_SERVICE_KIND=gateway
Environment=OPENCLAW_SERVICE_VERSION=2026.1.30

[Install]
WantedBy=default.target
